# config/spiderfoot.service
[Unit]
Description=SpiderFoot OSINT Engine
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/spiderfoot
ExecStart=/opt/spiderfoot/venv/bin/python3 /opt/spiderfoot/sf.py -l 127.0.0.1:5001
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
